# relay/smf/manifests/amon-relay.xml.in
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">
<!--
    This Source Code Form is subject to the terms of the Mozilla Public
    License, v. 2.0. If a copy of the MPL was not distributed with this
    file, You can obtain one at http://mozilla.org/MPL/2.0/.
-->

<!--
    Copyright (c) 2014, Joyent, Inc.
-->

<service_bundle type="manifest" name="amon-relay">
  <service name="smartdc/agent/amon-relay" type="service" version="@@VERSION@@">
    <create_default_instance enabled="true"/>
    <single_instance/>

    <dependency name="network" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/milestone/network:default"/>
    </dependency>
    <dependency name="filesystem" grouping="require_all" restart_on="error" type="service">
      <service_fmri value="svc:/system/filesystem/local"/>
    </dependency>

    <exec_method
      type="method"
      name="start"
      exec="/usr/bin/ctrun -l child -o noorphan @@PREFIX@@/lib/node_modules/amon-relay/bin/amon-relay -v -Z &amp;"
      timeout_seconds="10">
      <method_context>
        <method_credential user="root" group="staff"/>
        <method_environment>
          <envvar name="LD_PRELOAD_32" value="/usr/lib/extendedFILE.so.1" />
        </method_environment>
      </method_context>
    </exec_method>

    <exec_method type="method" name="refresh" exec=":kill" timeout_seconds="10">
      <method_context>
        <method_credential user="root" group="staff"/>
      </method_context>
    </exec_method>

    <exec_method type="method" name="stop" exec=":kill" timeout_seconds="10">
      <method_context>
        <method_credential user="root" group="staff"/>
      </method_context>
    </exec_method>

    <property_group name="startd" type="framework">
      <propval name="ignore_error" type="astring" value="core,signal"/>
    </property_group>
    <property_group name="application" type="application" />

    <stability value="Evolving"/>

    <template>
      <common_name>
        <loctext xml:lang="C">SDC Amon Relay</loctext>
      </common_name>
    </template>
  </service>
</service_bundle>
